# debian/saferdefiance-jumpbox-daemon.postinst
#!/bin/sh

set -e

# Create a user specific for jumpbox
if ! getent passwd djb >/dev/null; then
	adduser --quiet --system --no-create-home --home /var/run --shell /usr/sbin/nologin djb
fi

# Make sure this is ours
chown -R djb:nobody /var/cache/saferdefiance/djb/

if [ -x /usr/sbin/invoke-rc.d ]; then
	invoke-rc.d djb restart || true
else
	/etc/init.d/djb restart || true
fi

#DEBHELPER#
